# root/etc/uci-defaults/80_advanced
#!/bin/sh

# 清理旧版 LuCI 缓存
rm -f /tmp/luci-indexcache
rm -f /tmp/luci-modulecache/*

exit 0
